# examples/mqtt/bin/atlas_ph_pub2mqtt.service
[Unit]
Description=Read Atlas pH metric from i2c bus and publish it.  See atlas_ph_pub2edge.sh for details on where.

[Service]
WorkingDirectory=/opt/pwnybrau/pwnybrau_atlas/examples/mqtt/bin
ExecStart=/bin/bash /opt/pwnybrau/pwnybrau_atlas/examples/mqtt/bin/atlas_ph_pub2edge.sh

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target  
